# templates/etc/network/if-up.d/iproute2-001-routes-Debian.erb
#!/bin/sh
###
### File managed by Puppet
###

[ "$IFACE" != "lo" ] || exit 0

#
# Add routes
#
<% @routes.each do |value| -%>
/sbin/ip route add <%= value['network'] %> via <%= value['gateway'] %> table <%= value['table'] %>
<% end -%>
