# templates/etc/network/if-up.d/iproute2-002-rules-Debian.erb
#!/bin/sh
###
### File managed by Puppet
###
<% @rules.each do |value| -%>
<%if value.has_key?('dev') %>
[ "$IFACE" == "<%= value['dev'] %>" ] || exit 0
<% else %>
[ "$IFACE" != "lo" ] || exit 0
<% end %>
/sbin/ip <% if value['ipv6'] -%>-6<% end -%> rule add from <%= value['from'] %><% if value.has_key?('to') %> to <%= value['to'] %><% end %> table <%= value['table'] %><% if value.has_key?('priority') %> priority <%= value['priority'] %><% end %>
<% end -%>
