# templates/etc/network/if-up.d/iproute2-001-routes-Debian.erb
#!/bin/sh
###
### File managed by Puppet
###

[ "$IFACE" != "lo" ] || exit 0

#
# Add routes
#
<% @routes.each do |value| -%>
/sbin/ip route add <%= value['network'] %><%if value.has_key?('gateway') %> via <%= value['gateway'] %><% end %> <%if value.has_key?('dev') %>dev <%= value['dev'] %><% end %><%if value.has_key?('src') %> src <%= value['src'] %><% end %> table <%= value['table'] %>
<% end -%>
