# templates/etc/network/if-up.d/iproute2-001-routes-Debian.erb
#!/bin/sh
###
### File managed by Puppet
###
<% @routes.each do |value| -%>
<%if value.has_key?('dev') %>
if [ "$IFACE" = "<%= value['dev'] %>" ]; then
<% else %>
if [ "$IFACE" != "lo" ]; then
<% end %>
  /sbin/ip route add <%= value['network'] %> <%if value.has_key?('dev') %>dev <%= value['dev'] %><% end %><%if value.has_key?('src') %> src <%= value['src'] %><% end %> table <%= value['table'] %>
  <%if value.has_key?('gateway') %>
  /sbin/ip route add default via <%= value['gateway'] %> table <%= value['table'] %>
  <%end%>
fi
<% end -%>
