# templates/etc/network/if-up.d/iproute2-002-rules-Debian.erb
#!/bin/sh
###
### File managed by Puppet
###

[ "$IFACE" != "lo" ] || exit 0

#
# Add rules
#
<% @rules.each do |value| -%>
/sbin/ip rule add from <%= value['from'] %><% if value.has_key?('to') %> to <%= value['to'] %><% end %> table <%= value['table'] %><% if value.has_key?('priority') %> priority <%= value['priority'] %><% end %>
<% end -%>
